# com.github.roon-tuneshine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.github.roon-tuneshine</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-W</string>
        <string>ignore</string>
        <string>/Users/setu/portfolios/roon-tuneshine/main.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Write stdout/stderr to a log file for debugging -->
    <key>StandardOutPath</key>
    <string>/tmp/roon-tuneshine.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/roon-tuneshine.log</string>

    <!-- Give the network time to come up before starting -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
